# variant/target/var/lib/dpkg/info/nchttpd.prerm
#!/bin/sh

invoke-rc.d httpd stop
update-rc.d httpd remove
