# variant/target/var/lib/dpkg/info/nchttpd.list
/etc/init.d/httpd
/usr/bin/httpd.sh
/var/www/api.xml
/var/www/default.xsl
/var/www/style.css
