# variant/target/etc/init.d/mountkernfs
#! /bin/sh
### BEGIN INIT INFO
# Provides:          mountkernfs
# Required-Start:
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Mount kernel virtual file systems.
# Description:       Mount initial set of virtual filesystems the kernel
#                    provides and that are required by everything.
### END INIT INFO

mount /proc
mount -a

mkdir /run/lock

mount -n -t devtmpfs mdev /dev
mkdir /dev/pts /dev/shm
mount -t devpts devpts /dev/pts
chmod 1777 /dev/shm
